# supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=/usr/bin/mongod

[program:mongodb-rest]
command=/usr/local/bin/mongodb-rest
